# launchd/com.bukitoka.tmetric-helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bukitoka.tmetric-helper</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>tmetric-helper</string>
        <string>auto-keep-active</string>
        <string>--inactivity-timeout</string>
        <string>300</string>
        <string>--check-interval</string>
        <string>10</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <!-- Check every 30 minutes (1800 seconds) -->
    <!-- This ensures the app tries to start periodically -->
    <!-- If outside work hours, it will exit immediately -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Only restart if it crashes unexpectedly (non-zero exit) -->
    <!-- Normal exit (exit code 0) when outside work hours won't restart -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/tmetric-helper.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/tmetric-helper.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Volumes/Devs/tmetric-helper</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Throttle settings to prevent rapid restarts -->
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>
